# compaction/crontab
*/30 * * * * /usr/bin/python3 /path/to/compaction/compact_parquet.py
